# mw_cgclient/supervisord/supervisord.conf
[supervisord]
nodaemon=true

# Ejecuto el servicio de ssh
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd_error.log

[program:tcpdump]
command=tcpdump -i eth0 -w /app/registro/cgclient.pcap
autostart=true
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd_error.log

# [program:ACROSSclient]
# command=python3 ACROSSclient.py
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/sshd.log
# stderr_logfile=/var/log/sshd_error.log

# [program:ACROSSfile_transfer]
# command=python3 ACROSSfile_transfer.py
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/sshd.log
# stderr_logfile=/var/log/sshd_error.log

# [program:ACROSSconsuming_video]
# command=python3 ACROSSconsuming_video.py
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/sshd.log
# stderr_logfile=/var/log/sshd_error.log

# [program:ACROSSshorts]
# command=python3 ACROSSshorts.py
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/sshd.log
# stderr_logfile=/var/log/sshd_error.log


